# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/ha-opencode/run
#!/command/with-contenv bashio
# ==============================================================================
# HA OpenCode - Main Service
# ==============================================================================

bashio::log.info "=============================================="
bashio::log.info "  HA OpenCode for Home Assistant"
bashio::log.info "  Starting services..."
bashio::log.info "=============================================="

# Ensure data directories exist
mkdir -p /data/.local/share/opencode
mkdir -p /data/.config/opencode

# Set up persistent HOME for serve mode
export HOME="/data"
export XDG_DATA_HOME="/data/.local/share"
export XDG_CONFIG_HOME="/data/.config"

# ==============================================================================
# Deploy AGENTS.md to Home Assistant config directory (first install only)
# Users can customize this file using File Editor or any text editor
# ==============================================================================

if [ ! -f "/homeassistant/AGENTS.md" ]; then
    cp /opt/ha-mcp-server/AGENTS.md /homeassistant/AGENTS.md
    bashio::log.info "Deployed default AGENTS.md to /homeassistant/"
    bashio::log.info "You can customize AI instructions by editing AGENTS.md"
fi

# ==============================================================================
# Read configuration options
# ==============================================================================

MCP_ENABLED=$(bashio::config 'mcp_enabled')
LSP_ENABLED=$(bashio::config 'lsp_enabled' || echo "true")
SERVE_PRINT_LOGS=$(bashio::config 'serve_print_logs' || echo "true")
SERVE_LOG_LEVEL=$(bashio::config 'serve_log_level' || echo "DEBUG")
SERVE_HOSTNAME=$(bashio::config 'serve_hostname' || echo "0.0.0.0")
SERVE_MDNS=$(bashio::config 'serve_mdns' || echo "false")
OPENCODE_SERVER_PASSWORD=$(bashio::config 'opencode_server_password' || echo "")
OPTIONS_FILE="/data/options.json"
SERVE_PORT=4096

bashio::log.info "Configuration: MCP=${MCP_ENABLED} LSP=${LSP_ENABLED}"

if [ "${OPENCODE_SERVER_PASSWORD}" = "null" ]; then
    OPENCODE_SERVER_PASSWORD=""
fi

if [ -n "${OPENCODE_SERVER_PASSWORD}" ]; then
    export OPENCODE_SERVER_PASSWORD="${OPENCODE_SERVER_PASSWORD}"
else
    unset OPENCODE_SERVER_PASSWORD
fi

# ==============================================================================
# Apply configuration in single jq operation
# ==============================================================================

MCP_BOOL=$([ "${MCP_ENABLED}" = "true" ] && echo "true" || echo "false")
LSP_DISABLED=$([ "${LSP_ENABLED}" = "true" ] && echo "false" || echo "true")

jq --argjson mcp "${MCP_BOOL}" --argjson lsp_disabled "${LSP_DISABLED}" '
    .mcp.homeassistant.enabled = $mcp |
    if $lsp_disabled then .lsp["ha-yaml"].disabled = true else del(.lsp["ha-yaml"].disabled) end
' /opt/ha-mcp-server/opencode-ha.json > /data/.config/opencode/opencode.json

# ==============================================================================
# Inject OpenCode config.json (optional)
# ==============================================================================

if [ -f "${OPTIONS_FILE}" ]; then
    OPENCODE_CONFIG_JSON=$(jq -r '.opencode_config // empty' "${OPTIONS_FILE}")
    if [ -n "${OPENCODE_CONFIG_JSON}" ]; then
        if ! printf '%s' "${OPENCODE_CONFIG_JSON}" | jq -e . >/dev/null 2>&1; then
            bashio::log.error "Invalid JSON in opencode_config"
            exit 1
        fi

        printf '%s' "${OPENCODE_CONFIG_JSON}" > /data/.config/opencode/config.json
        bashio::log.info "Wrote /data/.config/opencode/config.json from add-on config"
    fi
fi

# ==============================================================================
# Set Node.js memory limit (hardcoded for stability)
# ==============================================================================

bashio::log.info "Starting OpenCode server on port ${SERVE_PORT}..."

# ==============================================================================
# Start OpenCode server
# ==============================================================================

SERVE_ARGS=(serve --hostname "${SERVE_HOSTNAME}" --port "${SERVE_PORT}" --log-level "${SERVE_LOG_LEVEL}")

if [ "${SERVE_PRINT_LOGS}" = "true" ]; then
    SERVE_ARGS+=(--print-logs)
fi

if [ "${SERVE_MDNS}" = "true" ]; then
    SERVE_ARGS+=(--mdns)
fi

if [ -f "${OPTIONS_FILE}" ]; then
    SERVE_CORS_TYPE=$(jq -r 'if .serve_cors==null then "null" else (.serve_cors|type) end' "${OPTIONS_FILE}")
    case "${SERVE_CORS_TYPE}" in
        array)
            mapfile -t SERVE_CORS_VALUES < <(jq -r '.serve_cors[] | select(length>0)' "${OPTIONS_FILE}")
            ;;
        string)
            mapfile -t SERVE_CORS_VALUES < <(jq -r '.serve_cors | split(",") | map(gsub("^\\s+|\\s+$"; "")) | map(select(length>0))[]' "${OPTIONS_FILE}")
            ;;
        null)
            SERVE_CORS_VALUES=()
            ;;
        *)
            bashio::log.error "Invalid serve_cors value; expected a list or comma-separated string"
            exit 1
            ;;
    esac

    for cors_value in "${SERVE_CORS_VALUES[@]}"; do
        SERVE_ARGS+=(--cors "${cors_value}")
    done
fi

exec opencode "${SERVE_ARGS[@]}"
