# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/ha-opencode/run
#!/command/with-contenv bashio
# ==============================================================================
# HA OpenCode - Main Service
# ==============================================================================

bashio::log.info "=============================================="
bashio::log.info "  HA OpenCode for Home Assistant"
bashio::log.info "  Starting services..."
bashio::log.info "=============================================="

# Ensure data directories exist
mkdir -p /data/.local/share/opencode
mkdir -p /data/.config/opencode

# ==============================================================================
# Deploy AGENTS.md to Home Assistant config directory (first install only)
# Users can customize this file using File Editor or any text editor
# ==============================================================================

if [ ! -f "/homeassistant/AGENTS.md" ]; then
    cp /opt/ha-mcp-server/AGENTS.md /homeassistant/AGENTS.md
    bashio::log.info "Deployed default AGENTS.md to /homeassistant/"
    bashio::log.info "You can customize AI instructions by editing AGENTS.md"
fi

# ==============================================================================
# Read configuration options
# ==============================================================================

MCP_ENABLED=$(bashio::config 'mcp_enabled')
LSP_ENABLED=$(bashio::config 'lsp_enabled' || echo "true")
bashio::log.info "Configuration: MCP=${MCP_ENABLED} LSP=${LSP_ENABLED}"

# ==============================================================================
# Apply configuration in single jq operation
# ==============================================================================

MCP_BOOL=$([ "${MCP_ENABLED}" = "true" ] && echo "true" || echo "false")
LSP_DISABLED=$([ "${LSP_ENABLED}" = "true" ] && echo "false" || echo "true")

jq --argjson mcp "${MCP_BOOL}" --argjson lsp_disabled "${LSP_DISABLED}" '
    .mcp.homeassistant.enabled = $mcp |
    if $lsp_disabled then .lsp["ha-yaml"].disabled = true else del(.lsp["ha-yaml"].disabled) end
' /opt/ha-mcp-server/opencode-ha.json > /data/.config/opencode/opencode.json

# ==============================================================================
# Set Node.js memory limit (hardcoded for stability)
# ==============================================================================

bashio::log.info "Starting OpenCode web UI on port 4096..."

# ==============================================================================
# Start OpenCode web UI (serve mode)
# ==============================================================================

exec opencode serve --host 0.0.0.0 --port 4096
