# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/ha-opencode/run
#!/command/with-contenv bashio
# ==============================================================================
# HA OpenCode - Main Service
# ==============================================================================

bashio::log.info "=============================================="
bashio::log.info "  HA OpenCode for Home Assistant"
bashio::log.info "  Starting services..."
bashio::log.info "=============================================="

# Ensure data directories exist
mkdir -p /data/.local/share/opencode
mkdir -p /data/.config/opencode

# ==============================================================================
# Deploy AGENTS.md to Home Assistant config directory (first install only)
# Users can customize this file using File Editor or any text editor
# ==============================================================================

if [ ! -f "/homeassistant/AGENTS.md" ]; then
    cp /opt/ha-mcp-server/AGENTS.md /homeassistant/AGENTS.md
    bashio::log.info "Deployed default AGENTS.md to /homeassistant/"
    bashio::log.info "You can customize AI instructions by editing AGENTS.md"
fi

# ==============================================================================
# Read configuration options
# ==============================================================================

MCP_ENABLED=$(bashio::config 'mcp_enabled')
LSP_ENABLED=$(bashio::config 'lsp_enabled' || echo "true")
OPENCODE_CONFIG_JSON=$(bashio::config 'opencode_config_json' || echo "")
SERVE_PRINT_LOGS=$(bashio::config 'serve_print_logs' || echo "true")
SERVE_LOG_LEVEL=$(bashio::config 'serve_log_level' || echo "DEBUG")
SERVE_PORT=$(bashio::config 'serve_port' || echo "4096")
SERVE_HOSTNAME=$(bashio::config 'serve_hostname' || echo "0.0.0.0")
SERVE_MDNS=$(bashio::config 'serve_mdns' || echo "false")
SERVE_CORS_RAW=$(bashio::config 'serve_cors' || echo "[]")

bashio::log.info "Configuration: MCP=${MCP_ENABLED} LSP=${LSP_ENABLED}"

# ==============================================================================
# Apply configuration in single jq operation
# ==============================================================================

MCP_BOOL=$([ "${MCP_ENABLED}" = "true" ] && echo "true" || echo "false")
LSP_DISABLED=$([ "${LSP_ENABLED}" = "true" ] && echo "false" || echo "true")

jq --argjson mcp "${MCP_BOOL}" --argjson lsp_disabled "${LSP_DISABLED}" '
    .mcp.homeassistant.enabled = $mcp |
    if $lsp_disabled then .lsp["ha-yaml"].disabled = true else del(.lsp["ha-yaml"].disabled) end
' /opt/ha-mcp-server/opencode-ha.json > /data/.config/opencode/opencode.json

# ==============================================================================
# Inject OpenCode config.json (optional)
# ==============================================================================

if [ -n "${OPENCODE_CONFIG_JSON}" ] && [ "${OPENCODE_CONFIG_JSON}" != "null" ]; then
    if ! printf '%s' "${OPENCODE_CONFIG_JSON}" | jq -e . >/dev/null 2>&1; then
        bashio::log.error "Invalid JSON in opencode_config_json"
        exit 1
    fi

    printf '%s' "${OPENCODE_CONFIG_JSON}" > /data/.config/opencode/config.json
    bashio::log.info "Wrote /data/.config/opencode/config.json from add-on config"
fi

# ==============================================================================
# Set Node.js memory limit (hardcoded for stability)
# ==============================================================================

bashio::log.info "Starting OpenCode web UI on port ${SERVE_PORT}..."

# ==============================================================================
# Start OpenCode web UI (serve mode)
# ==============================================================================

SERVE_ARGS=(serve --hostname "${SERVE_HOSTNAME}" --port "${SERVE_PORT}" --log-level "${SERVE_LOG_LEVEL}")

if [ "${SERVE_PRINT_LOGS}" = "true" ]; then
    SERVE_ARGS+=(--print-logs)
fi

if [ "${SERVE_MDNS}" = "true" ]; then
    SERVE_ARGS+=(--mdns)
fi

if printf '%s' "${SERVE_CORS_RAW}" | jq -e . >/dev/null 2>&1; then
    mapfile -t SERVE_CORS_VALUES < <(printf '%s' "${SERVE_CORS_RAW}" | jq -r '.[]')
    for cors_value in "${SERVE_CORS_VALUES[@]}"; do
        SERVE_ARGS+=(--cors "${cors_value}")
    done
else
    bashio::log.error "Invalid serve_cors value; expected a list of strings"
    exit 1
fi

exec opencode "${SERVE_ARGS[@]}"
